# debian/iwd.postinst
# enable sysvinit script
update-rc.d iwd defaults

#DEBHELPER#
