# debian/iwd.postrm
#!/bin/sh
set -e

# remove sysvinit script
update-rc.d iwd remove

#DEBHELPER#
